# pantheon-core/pantheon-mcp.service
[Unit]
Description=Pantheon MCP Server — exposes Athenaeum, messaging, and god systems as MCP tools
Documentation=https://github.com/Duskript/Pantheon-Core
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/konan/pantheon/pantheon-core
ExecStart=/home/konan/.hermes/hermes-agent/venv/bin/python3 /home/konan/pantheon/pantheon-core/mcp_server.py --port 8010
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
