# scripts/clawforge-federation-stats.service
[Unit]
Description=Clawforge Federation Stats (one-shot)
After=network-online.target clawforge-registry.service nats.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /usr/local/bin/clawforge-federation-stats.py
User=root
Group=root
StandardOutput=journal
StandardError=journal
# Federation stats runs as a one-shot triggered by the matching .timer.
# No Restart= — failures are visible in journalctl.
